# packages/bitcore-node/bitcore-node.service
[Unit]
Description=Bitcore Node
After=network.target

[Service]
User=stratis
Group=stratis
Environment="BITCORE_CONFIG_PATH=/home/stratis/bitcore/bitcore.config.json"
ExecStart=/usr/bin/node /home/stratis/Desktop/bitcore/packages/bitcore-node/build/src/server.js

[Install]
WantedBy=multi-user.target